# snapraid/etc/systemd/system/snapraid-scrub.timer
[Unit]
Description=Snapraid Scrub Timer

[Timer]
Unit=snapraid-scrub.service
OnCalendar=Mon *-*-* 02:00:00

[Install]
WantedBy=timers.target
