# snapraid/config/etc/systemd/system/snapraid-scrub.service
[Unit]
Description=Snapraid Scrub Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/snapraid-scrub.sh
