# snapraid/config/etc/systemd/system/snapraid-sync.service
[Unit]
Description=SnapRAID Sync Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/snapraid-sync.sh
