# snapraid/config/etc/systemd/system/storage-maintenance.service
[Unit]
Description=Storage Maintenance

[Service]
Type=oneshot
ExecStart=/opt/scripts/storage-maintenance-unattended.sh
